# supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=/elasticsearch/bin/elasticsearch
autorestart=true

[program:kibana]
command=/kibana/bin/kibana
autorestart=true
